# dev_sandbox_starter_files/cmp.ini
[program:laravel-worker-group]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=default --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=5
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/default-worker.log

[program:laravel-worker-fileDownloads]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=fileDownloads --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/file-worker.log

[program:laravel-worker-BlueHornet]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=BlueHornet --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=5
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/blueHornet-worker.log

[program:Daddy-Warbucks]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=orphanage --tries=2
autostart=true
autorestart=true
user=cmpadmin
startretries=2
startsecs = 0
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/orphan-worker.log

[program:Campaigner]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=Campaigner --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/campaigner-worker.log

[program:Publicators]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=Publicators --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/publicators-worker.log

[program:attribution]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=attribution --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=5
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/attribution-worker.log

[program:filters]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=filters --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/filters-worker.log

[program:AWeber]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=AWeber --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/aweber-worker.log

[program:Maro]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=Maro --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/maro-worker.log

[program:Bronto]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=Bronto --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/bronto-worker.log

[program:Monitor]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=Monitor --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=1
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/monitor-worker.log

[program:Notifications]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=scheduled_notifications --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=1
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/notification-worker.log

[program:rawFeedProcessing]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=rawFeedProcessing --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/rawFeedProcessing-worker.log

[program:RecordProcessing]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/cmp/artisan queue:work --queue=RecordProcessing --tries=2
autostart=true
autorestart=true
user=cmpadmin
numprocs=2
startretries=2
startsecs = 0
redirect_stderr=true
stdout_logfile=/var/www/html/cmp/storage/logs/recordprocessing-worker.log
